# target_files/obmc-bmc-passive.target
[Unit]
Description=Passive BMC Target
Conflicts=obmc-bmc-active.target
