# target_files/obmc-bmc-active.target
[Unit]
Description=Active BMC Target
Conflicts=obmc-bmc-passive.target
